# rootfs/etc/update-motd.d/11-image-support
#!/bin/sh
#
#    11-help-text - print the help text associated with image

printf " * Image Support:  https://github.com/mengzhuo/hikey970-ubuntu-image\n"
